# elements/app_fum/install.d/systeminfo.service
[Unit]
Description=System Info

[Service]
ExecStart=/bin/bash /root/codes/start_sys.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
